# systemd/user/xdg-desktop-portal-wlr.service.d/override.conf
[Service]
Environment=SELECT_SOURCE=tofi
